# .config/systemd/user/awesome.service
[Unit]
Description=awesome
After=xorg.target

[Service]
ExecStart=/usr/bin/awesome
Environment=DISPLAY=:0

[Install]
WantedBy=awesome.target
