# .config/systemd/user/urxvtd.service
[Unit]
Description=Urxvt Daemon
After=xorg.target

[Service]
Type=simple
ExecStart=/usr/bin/urxvtd

[Install]
WantedBy=xinit.target
