# .config/systemd/user/xinit.target
[Unit]
Description=Xinit
Requires=xorg.target
