# .config/systemd/user/xmonad.service
[Unit]
Description=xmonad
After=xorg.target

[Service]
ExecStart=/usr/bin/xmonad
Environment=DISPLAY=:0

[Install]
WantedBy=xmonad.target
